# systemd/power-stats.timer
[Unit]
Description=Run power-stats.service every minute

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=timers.target
